# imooc_paid.conf
[program:imooc_paid]
command=python3 run.py
directory=/home/crawler/imooc_paid
autostart=true
autorestart=true
environment=LANG="en_US.utf8", LC_ALL="en_US.UTF-8", LC_LANG="en_US.UTF-8"